# jellywatchd.service
[Unit]
Description=JellyWatch Media Organizer Daemon
Documentation=https://github.com/Nomadcxx/jellywatch
After=network.target

[Service]
Type=simple
# Change to your user/group
User=jellyfin
Group=jellyfin
ExecStart=/usr/local/bin/jellywatchd
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# Add paths to your media directories here
# ReadWritePaths=/path/to/media /path/to/downloads
PrivateTmp=true

[Install]
WantedBy=multi-user.target
