# jellywatchd.service
[Unit]
Description=JellyWatch Media Organizer Daemon
Documentation=https://github.com/Nomadcxx/jellywatch
After=network.target

[Service]
Type=simple
User=jellyfin
Group=jellyfin
ExecStart=/usr/local/bin/jellywatchd --health-addr :8686
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Health check for systemd
ExecStartPost=/usr/bin/timeout 30 sh -c 'until curl -sf http://localhost:8686/ready; do sleep 1; done'

# Graceful shutdown (SIGTERM, then SIGKILL after timeout)
TimeoutStopSec=30
KillMode=mixed

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

# Uncomment and customize these paths for your setup:
# ReadWritePaths=/path/to/media /path/to/downloads

[Install]
WantedBy=multi-user.target
